# files/redhat/usr/share/cassandra/etc/cron.d/cassandra.cron
# Run watchdog.sh every min
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root
HOME=/

* * * * * root /usr/share/cassandra/sbin/watchdog.sh > /var/log/cassandra/watchdog.log
